# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:slskd]
command=/opt/slskd/slskd --app-dir=/usr/src/app/slskd-data
autostart=true
autorestart=true
startsecs=3
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=SLSKD_NO_AUTH="false",SLSKD_API_KEY="%(ENV_SLSKD_API_KEY)s",SLSKD_SLSK_USERNAME="%(ENV_SLSK_USERNAME)s",SLSKD_SLSK_PASSWORD="%(ENV_SLSK_PASSWORD)s",SLSKD_DOWNLOADS_DIR="/downloads",HOME="/root"

[program:bot]
command=npm start
directory=/usr/src/app
autostart=true
autorestart=true
startsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
